# deploy/metricshour-api.service
[Unit]
Description=MetricsHour FastAPI
After=network.target

[Service]
User=root
WorkingDirectory=/root/metricshour/backend
Environment="PATH=/root/metricshour/workers/venv/bin"
EnvironmentFile=/root/metricshour/backend/.env
ExecStart=/root/metricshour/workers/venv/bin/gunicorn app.main:app \
    -w 2 \
    -k uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --access-logfile /var/log/metricshour/api-access.log \
    --error-logfile /var/log/metricshour/api-error.log \
    --timeout 60
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
